# lowbat_hibernate/.config/systemd/user/lowbat_hibernate.service
[Unit]
Description=Run script to hibernate if low bat

[Service]
Environment=DISPLAY=:0
Type=oneshot
ExecStart=/home/jrutledge/bin/lowbat-hibernate
